# evaluation/fastfreeze/s6-rc.d/minio/run
#!/command/with-contenv sh
exec /usr/sbin/minio server /data