# source/pind/s6-rc.d/interface/run
#!/command/with-contenv sh
cd /app
exec npm run start
