# source/nested-container-runtime/s6-rc.d/dind/run
#!/command/with-contenv sh
exec dockerd-entrypoint.sh dockerd --tls=false --host=tcp://0.0.0.0:2375